# etc/udev/rules.d/99-garmin-strava.rules
ATTRS{idVendor}=="091e", ATTRS{idProduct}=="25ca", RUN+="/home/pi/strava-udev/scripts/udev.sh"
